# root_laptop/etc/udev/rules.d/70-monitor.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="autorandr --change --default laptop"
